# roles/initweb/files/initweb.service
[Unit]
Description=Initweb service

[Service]
ExecStart=/home/lvi/initweb.sh

[Install]
WantedBy=multi-user.target